# etc/init.d/sqtd
#!/sbin/runscript
# Copyright Vladimir Vershinin, Russia, Volgograd
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-analyzer/sqtd/files/2.0/init.d/sqtd,v 0.0.0 2013/03/08 10:53:34 mattm Exp $

pid_file="/var/lib/sqtd/sqtd.pid"

depend() {
	use logger
	provide sqtd
}

start_pre() {
	checkpath -d -m 0775 -o root:root /var/lib/sqtd
}

start() {
	ebegin "Starting sqtd"
	start-stop-daemon --start --wait 1000 \
		--user root --group root --exec /usr/sbin/sqtd -- -c /etc/sqtd/sqtd.conf ${SQTD_DEBUGLEVEL}
	eend $?
}

stop() {
	ebegin "Stopping sqtd"
	start-stop-daemon  --stop --pidfile ${pid_file}
	eend $?
}
